# debianservices/myservice.service
[Unit]
Description=My Service
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/local/bin/myservice

[Install]
WantedBy=multi-user.target
